# containers/fcron/supervisord.conf
[supervisord]
nodaemon=true

[program:fcron]
command=/usr/local/sbin/fcron -f
stdout_logfile=/var/log/supervisor/supervisord.log
stderr_logfile=/var/log/supervisor/supervisord.log
autorestart=true
